# etc/cron.d/affa
*/15 * * * * root /sbin/affa --_cronupdate
